# launchagents/com.keke.portal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.keke.portal</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/vvhome/.openclaw/workspace/projects/keke-vocab/apps/portal/scripts/start-daemon.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/vvhome/.openclaw/workspace/projects/keke-vocab/apps/portal</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/vvhome/.openclaw/workspace/projects/keke-vocab/apps/portal/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vvhome/.openclaw/workspace/projects/keke-vocab/apps/portal/logs/stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>ProcessType</key>
    <string>Background</string>
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
